# distros/alprd.service
[Unit]
Description=Automatic License Plate Recognition library

[Service]
ExecStart=/usr/bin/alprd -f

[Install]
WantedBy=multi-user.target
